# etc/cron.d/gs-carddav-reader-to-db
# Gemeinschaft
# $Revision$

# refresh cloud card phone books into the database

*/17 *  * * *   root   sleep 23; /opt/gemeinschaft/sbin/gs-carddav-reader.php 1>>/dev/null 2>>/dev/null
