# etc/cron.d/gs-cc-guardian
# Gemeinschaft
# $Revision$

# call completion

* *  * * *   root   sleep 2; /opt/gemeinschaft/sbin/gs-cc-guardian 1>>/dev/null 2>>/dev/null

