# etc/cron.d/gs-prov-checker
# Gemeinschaft
# $Revision$

# execute provisioning jobs

* *  * * *   root   sleep 4; /opt/gemeinschaft/sbin/gs-prov-checker 1>>/dev/null 2>>/dev/null

