# etc/cron.d/gs-prov-checker
# Gemeinschaft
# $Revision: 3261 $

# do open prov jobs

* *  * * *   root   sleep 8; /opt/gemeinschaft/sbin/gs-prov-checker 1>>/dev/null 2>>/dev/null

